# lib/code-server/model/supervisord-conf/include/code-server.conf
[program:code-server]
command=%(ENV_CODE_SERVER_DIR)s/bin/code-server --bind-addr=%(ENV_BIND_ADDR)s --auth password --user-data-dir=%(ENV_USER_DATA_PATH)s --config %(ENV_CONFIG_PATH)s  --disable-update-check --locale=zh-cn --extensions-dir=%(ENV_EXTENSION_PATH)s
autostart=false
autorestart=true
stopwaitsecs=3
stdout_logfile=%(ENV_CODE_SERVER_RUN_DIR)s/supervisord-conf/code-server.log
redirect_stderr=true
